# etc/systemd/system/kubelet.service.d/01-ibm-kubelet-addition.conf
[Service]
EnvironmentFile=/etc/sysconfig/ibmmachinemetadataenvfile
