# etc/systemd/system/crio.service.d/01-ibm-tmpdir.conf
[Service]
Environment="TMPDIR=/var/data/tmp"
